# fabzclean.service
[Unit]
Description=FabZClean Web Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/fabzclean
Environment="NODE_ENV=production"
EnvironmentFile=/opt/fabzclean/.env
ExecStart=/usr/bin/node /opt/fabzclean/start.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fabzclean

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (adjust based on your server capacity)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

