# systemd/fabzclean.service
[Unit]
Description=FabZClean Gunicorn Service
After=network.target mysql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/fabzclean
Environment="FLASK_APP=src.app"
Environment="FLASK_ENV=production"
Environment="DATABASE_URL=mysql+pymysql://fabz:fabzpass@127.0.0.1:3306/fabzclean"
Environment="SECRET_KEY=change-me-in-production"
Environment="JWT_SECRET_KEY=change-me-in-production"
Environment="ALLOW_ORIGINS=https://your-domain.com,file://"

ExecStart=/usr/local/bin/gunicorn \
    -w 3 \
    -b unix:/run/fabzclean.sock \
    --timeout 120 \
    --access-logfile /var/log/fabzclean/access.log \
    --error-logfile /var/log/fabzclean/error.log \
    "src.app:create_app('prod')"

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

